# conf/supervisord/gunicorn.conf
[program:gunicorn]
command=/app/bin/run-gunicorn
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
