# conf/supervisord/nginx.conf
[program:nginx]
command=/app/bin/run-nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
